# babybuddy/root/etc/services.d/babybuddy/run
#!/command/with-contenv bashio
#

declare -a options

options+=(-b :8000)
options+=(--log-level="$(bashio::config 'log_level')")
options+=(--worker-tmp-dir=/dev/shm)
options+=(--log-file=-)
options+=(--workers=2)
options+=(--threads=4)
options+=(--worker-class=gthread)

cd /app/babybuddy || bashio.exit.nok 'Failed cd, /app/babybuddy does not exist, exiting...'

export \
    DJANGO_SETTINGS_MODULE="babybuddy.settings.base" \
    ALLOWED_HOSTS="${ALLOWED_HOSTS:-*}" \
    TIME_ZONE="${TZ:-UTC}" \
    DEBUG="${DEBUG:-False}" \
    SECRET_KEY="${SECRET_KEY:-`cat /config/.secretkey`}"

export \
    CSRF_TRUSTED_ORIGINS=$(bashio::config "CSRF_TRUSTED_ORIGINS") \
    NAP_START_MIN=$(bashio::config "NAP_START_MIN") \
    NAP_START_MAX=$(bashio::config "NAP_START_MAX") \
    USE_24_HOUR_TIME_FORMAT=$(bashio::config "USE_24_HOUR_TIME_FORMAT")

if bashio::config.has_value "INGRESS_USER"; then
    bashio::log.info 'Adding config for Ingress User Auth'
    export \
        REVERSE_PROXY_AUTH=true \
        PROXY_HEADER=HTTP_X_REMOTE_USER
fi

if bashio::config.has_value "DB_HOST"; then
    bashio::log.info 'Running with user defined database'
    export \
        DB_ENGINE=$(bashio::config "DB_ENGINE") \
        DB_HOST=$(bashio::config "DB_HOST") \
        DB_NAME=$(bashio::config "DB_NAME") \
        DB_PASSWORD=$(bashio::config "DB_PASSWORD") \
        DB_PORT=$(bashio::config "DB_PORT") \
        DB_USER=$(bashio::config "DB_USER")
fi

exec \
    s6-setuidgid root gunicorn babybuddy.wsgi "${options[@]}"

